# if-up.d/ip6tables
#!/bin/sh
/sbin/ip6tables-restore < /etc/ip6tables.rules